# container/s6-overlay/s6-rc.d/movim-migrations/up
with-contenv composer -d /var/www/movim movim:migrate
