# container/s6-overlay/s6-rc.d/movim-migrations/up
with-contenv
s6-setuidgid www-data composer -d /var/www/movim movim:migrate
