# container/s6-overlay/s6-rc.d/tail-log/run
#!/command/with-contenv sh
tail -F /var/www/movim/log/errors.log /var/www/movim/log/debug.log /var/www/movim/log/info.log -q 2> /dev/null
